# distask.service
[Unit]
Description=DisTask Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/distask
EnvironmentFile=%h/distask/.env
ExecStart=%h/distask/.venv/bin/python bot.py
Restart=always
RestartSec=5
User=%i

[Install]
WantedBy=multi-user.target
