# distask-web.service
[Unit]
Description=DisTask Landing Page
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/distask
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/distask/.venv/bin/uvicorn web.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=%i

[Install]
WantedBy=multi-user.target
