# distask-git-sync.timer
[Unit]
Description=DisTask GitHub Sync Timer (Fallback Polling)

[Timer]
OnCalendar=*:0/2
Persistent=true
RandomizedDelaySec=30
Unit=distask-git-sync.service

[Install]
WantedBy=timers.target

